# deploy/com.example.qa-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  web-qa-agent — haftalık QA taraması (launchd LaunchAgent, ÖRNEK).

  Her Pazartesi 09:00'da weekly_run.py'yi (venv python ile) çalıştırır.

  ‼️ ÖNCE YOLLARI DÜZENLEYİN: aşağıdaki /path/to/web-qa-agent ifadelerini kendi
     kurulum dizininizle değiştirin (ör. /Users/<kullanıcı>/web-qa-agent).

  Kurulum:
    cp deploy/com.example.qa-weekly.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.example.qa-weekly.plist

  Elle test (hemen bir kez çalıştır):
    launchctl start com.example.qa-weekly
    # veya modern:
    launchctl kickstart -k gui/$(id -u)/com.example.qa-weekly

  Durdurma/kaldırma:
    launchctl unload -w ~/Library/LaunchAgents/com.example.qa-weekly.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.qa-weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/path/to/web-qa-agent/venv/bin/python</string>
        <string>/path/to/web-qa-agent/weekly_run.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/path/to/web-qa-agent</string>

    <!-- Her Pazartesi 09:00 (Weekday: 0/7=Pazar, 1=Pazartesi ... 6=Cumartesi) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- launchd job'un çıktısı (weekly_run zaten logs/weekly.log'a ayrıca yazar) -->
    <key>StandardOutPath</key>
    <string>/path/to/web-qa-agent/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/web-qa-agent/logs/launchd.err.log</string>

    <!-- Yüklerken hemen çalışmasın; sadece takvim tetiklesin (elle test edeceğiz). -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
